# supervisord.conf
[supervisord]
nodaemon=true

[program:haproxy]
command=/usr/local/sbin/haproxy -db -f /usr/local/etc/haproxy/haproxy.cfg

[program:marathon-check]
command=/usr/bin/marathon-check.sh
